# systemd.timer
[Unit]
Description=Run vivaldi-profile-backup every few hours

[Timer]
OnBootSec=10m
OnActiveSec=2h

[Install]
WantedBy=timers.target